# packer/rhel7/config/etc/profile.d/harden.sh
umask 027
readonly TMOUT=900
HISTCONTROL=ignoredups
readonly HISTCONTROL HISTFILE HISTIGNORE HISTFILESIZE
mesg n 

