# packer/apigw/config/rbs-ecomm-db-purge.cron
# Purge all non-severe audit records over 3 weeks old
0 0 * * * root /usr/local/bin/audit_purge.sh -p
